# templates/default/oracle.erb
#!/bin/bash
#
# oracle
#
# chkconfig: 2345 99 01
# description:	Starts, stops oracle databases
#
# Source function library.

#
# !!!!! THIS FILE GENERATED BY CHEF  !!!!!
# !!!!! ANY CHANGES WILL BE REVERTED !!!!!
#


. /etc/init.d/functions


if [ ${EUID} -gt 0 ] ; then
   echo "run as root."
   exit 1
fi

start() {
    touch /var/lock/subsys/oracle
    if [ -f /etc/oratab ] ; then
       while  read -r line  ; do
          ch=`echo ${line} | cut -c1-1`
          if [ "${ch}" == "#" -o "${ch}" == "" ] ; then
             continue
          fi
          auto=`echo ${line} | awk -F: '{print $3}'`
          export ORACLE_SID=`echo ${line} | awk -F: '{print $1}'`
          export ORACLE_HOME=`echo ${line} | awk -F: '{print $2}'`
          export PATH=/usr/lib/jvm/default/bin:/usr/local/sbin:/usr/bin:/bin:/usr/sbin:/sbin:/root/bin:${ORACLE_HOME}/bin
          if [ ${auto} == "Y" ] ; then
             running=`ps -ef | grep " ora_pmon_${ORACLE_SID}\$" | grep -v grep | wc -l`
             if [ $running -eq 0 ] ; then
                su oracle -c 'sqlplus "/ as sysdba" << EOT
                    startup
                    exit
EOT
    lsnrctl start <%= @oracle_sid.upcase %>LISTENER
    '
             fi
          fi
       done < /etc/oratab
    
    fi
    
    return 0
}

stop() {

   rm -f /var/lock/subsys/oracle
   if [ -f /etc/oratab ] ; then
      while  read -r line  ; do
      echo ${line}
         ch=`echo ${line} | cut -c1-1`
         if [ "${ch}" == "#" -o "${ch}" == "" ] ; then
            continue
         fi
         export ORACLE_SID=`echo ${line} | awk -F: '{print $1}'`
         export ORACLE_HOME=`echo ${line} | awk -F: '{print $2}'`
         export PATH=/usr/lib/jvm/default/bin:/usr/local/sbin:/usr/bin:/bin:/usr/sbin:/sbin:/root/bin:${ORACLE_HOME}/bin
   
         running=`ps -ef | grep " ora_pmon_${ORACLE_SID}\$" | grep -v grep | wc -l`
   
         if [ $running -gt 0 ] ; then
            su oracle -c 'sqlplus "/ as sysdba" << EOT
                shutdown immediate
                exit
            EOT
            lsnrctl stop <%= @oracle_sid.upcase %>LISTENER
       '
             fi
          done < /etc/oratab
       
       fi

        return 0
    }

case "$1" in
    start)
	start
	RETVAL=$?
	;;
    stop)
	stop
	RETVAL=$?
	;;
    *)
	RETVAL=2
	;;
esac

exit $RETVAL


